# contrib/dunst.systemd.service.in
[Unit]
Description=Dunst notification daemon

[Service]
Type=dbus
BusName=org.freedesktop.Notifications
ExecStart=##PREFIX##/bin/dunst
Environment=DISPLAY=:0

[Install]
WantedBy=default.target

