# dunst.systemd.service.in
[Unit]
Description=Dunst notification daemon
Documentation=man:dunst(1)
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=dbus
BusName=org.freedesktop.Notifications
ExecStart=@bindir@/dunst
ExecReload=@bindir@/dunstctl reload
Slice=session.slice
